# .config/startup/firstboot.service
[Unit]
Description=Post-install configuration run on first boot
[Service]
Type=simple
ExecStart=$HOME/.config/startup/post_install.sh
[Install]
WantedBy=multi-user.target